# systemd/barksignal-update.service
[Unit]
Description=BarkSignal Update (git pull)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/barksignal-update.sh
