# systemd/barksignal-portal.service
[Unit]
Description=BarkSignal Setup Portal
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=barksignal
WorkingDirectory=/home/barksignal/barksignal/portal
ExecStart=/home/barksignal/venv-barksignal/bin/gunicorn --workers 1 --threads 4 --bind 0.0.0.0:8080 --access-logfile - --error-logfile - app:app
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
